# shop-new.service
[Unit]
Description=Shop Fullstack (Backend + Frontend Dev)
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/shadcn-ui

# Запуск обоих процессов через bash
ExecStart=/bin/bash -c '\
  cd /home/ubuntu/shadcn-ui/backend && pnpm start & \
  cd /home/ubuntu/shadcn-ui && pnpm dev \
'

Restart=always
RestartSec=5
Environment=PATH=/usr/bin:/usr/local/bin:/home/ubuntu/.nvm/versions/node/v20.11.1/bin
Environment=PNPM_HOME=/home/ubuntu/.local/share/pnpm

[Install]
WantedBy=multi-user.target